# dyndns/crontabs/root.example
*/5 * * * * curl https://ddns_username:ddns_password@dyndns.kasserver.com 

